# deploy/d-brain-market.service
[Unit]
Description=d-brain Morning Market Digest
After=network.target

[Service]
Type=oneshot
User=shima
WorkingDirectory=/home/shima/projects/agent-second-brain
ExecStart=/home/shima/.local/bin/uv run python scripts/market.py
EnvironmentFile=/home/shima/projects/agent-second-brain/.env
StandardOutput=journal
StandardError=journal
